# backup.timer
[Unit]
Description=Run backup using Restic 5 minutes after boot

[Timer]
OnBootSec=5min
Unit=backup.service

[Install]
WantedBy=timers.target